# puppet/modules/dirvish/templates/systemd.dirvish.service.erb
[Unit]
Description=Run Dirvish backup jobs

[Service]
Type=oneshot
ExecStart=/etc/dirvish/dirvish-cronjob
User=root
Group=root
