# puppet/modules/redmine/templates/redmine.service.erb
[Unit]
Description=redmine

[Service]
Type=notify
User=<%= @username %>
ExecStart=<%= @app_root %>/bin/rails server --environment production
PrivateTmp=true
WorkingDirectory=<%= @app_root %>
SyslogIdentifier=redmine
Environment=RAILS_LOG_TO_STDOUT=1
