# puppet/modules/redmine/files/redmine.socket
# Deployed with puppet
#

[Unit]
Description=redmine socket

[Socket]
ListenStream=127.0.0.1:3000

[Install]
WantedBy=sockets.target
