# puppet/modules/koji/templates/mirror_reposync.cron.epp
<%- | Stdlib::Absolutepath $mirror_root | -%>
#!/bin/bash

dnf reposync \
  --config /etc/reposync.conf \
  --download-path=<%= $mirror_root %> \
  --download-metadata \
  --newest-only \
  --delete \
  --quiet
